# build/rootfs/etc/supervisor/conf.d/wine-vnc.conf
[group:wine-vnc]
programs=xvfb,wm,lxpanel,pcmanfm,x11vnc

[program:xvfb]
priority=10
command=/usr/local/bin/xvfb.sh
stopsignal=KILL

[program:wm]
priority=15
command=/usr/bin/openbox
environment=DISPLAY=":1",HOME="/root",USER="root"

[program:lxpanel]
priority=15
directory=/home/nobody
command=/usr/bin/lxpanel --profile LXDE
user=nobody
environment=DISPLAY=":1",HOME="/home/nobody",USER="nobody"

[program:pcmanfm]
priority=15
directory=/home/nobody
command=/usr/bin/pcmanfm --desktop --profile LXDE
user=nobody
environment=DISPLAY=":1",HOME="/home/nobody",USER="nobody"

[program:x11vnc]
priority=20
environment = DISPLAY=":1",XAUTHLOCALHOSTNAME="localhost"
command=bash -c "sleep 5 && x0vncserver -rfbauth .password2"

